# consul-template/consul-template.timer
[Unit]
Description=Run consul-template.service every minute

[Timer]
OnCalendar=*:0/1
